# i2c_screen_daemon.rc
#!/bin/bash
### BEGIN INIT INFO
# Provides:          i2c_screen_daemon
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Run i2c screen daemon
### END INIT INFO

DAEMON_PATH="/usr/local/bin/i2c_screen_daemon.py"
PIDFILE="/var/run/i2c_screen_daemon.pid"
LOGFILE="/var/log/i2c_screen_daemon.log"

start() {
  echo "Starting i2c_screen_daemon"
  if [ -f "$PIDFILE" ]; then
    PID=$(cat "$PIDFILE")
    if ps -p "$PID" > /dev/null 2>&1; then
      echo "Daemon is already running (PID $PID)"
      exit 0
    else
      echo "Stale PID file found. Removing."
      rm -f "$PIDFILE"
    fi
  fi

  nohup "$DAEMON_PATH" >> "$LOGFILE" 2>&1 &
  echo $! > "$PIDFILE"
  echo "Daemon started with PID $(cat $PIDFILE)"
}

stop() {
  echo "Stopping i2c_screen_daemon"
  if [ -f "$PIDFILE" ]; then
    PID=$(cat "$PIDFILE")
    if ps -p "$PID" > /dev/null 2>&1; then
      kill "$PID" && echo "Daemon stopped"
    else
      echo "PID file exists, but process is not running"
    fi
    rm -f "$PIDFILE"
  else
    echo "Daemon is not running"
  fi
}

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart)
    stop
    sleep 1
    start
    ;;
  *)
    echo "Usage: /etc/init.d/i2c_screen_daemon {start|stop|restart}"
    exit 1
    ;;
esac

exit 0
